# board/stp/raspberrypi0w/overlay/var/spool/cron/crontabs/root
* * * * * /usr/local/bin/reBOOT24.sh
1 * * * * /usr/local/bin/wifiPriority.sh
